# rpi-headless-config/rootfs/etc/crontab
@reboot     root    /boot/tasks.sh boot
